# provisioning/roles/redis/files/redis.conf
# redis
# description "An Upstart task to start redis

start on startup
setuid www-data
console log

exec /usr/local/bin/redis-server
