# provisioning/roles/redis/files/redis.conf
# simple redis script

description "redis server instance"
start on runlevel [2345]
stop on runlevel [06]
setuid www-data

respawn

chdir /var/www/app/

exec redis-server
